# app-misc/glsa-check/files/glsa-check.timer
[Timer]
OnCalendar=daily
AccuracySec=1m
Persistent=true

[Install]
WantedBy=multi-user.target
